# fk-es-cosmos-plugin/etc/cron.d/fk-es-cosmos-plugin
* * * * * root python /usr/share/es-cosmos-plugin/es-cosmos-plugin.py 2>/tmp/es-cosmos-plugin-errors.log | /usr/bin/cosmos